# systemd/system/xfs-internal.socket
[Unit]
Description=X Font Service (internal)

[Socket]
ListenStream=/run/xfs.sock
Service=xfs.service

[Install]
WantedBy=sockets.target
